# src/pangalink/build/modules/pangalink/files/etc/supervisor/conf.d/pangalink.conf
[program:pangalink]
command = /opt/pangalink-net/bin/pangalink-net --node-env="production" --config-dir="/etc/pangalink-net.d" --pid-file="/var/run/pangalink-net.pid"
autostart = true
autorestart = true
stdout_events_enabled = true
stderr_events_enabled = true
